# host-service/copy-passwd.path
[Unit]
Description=Copy passwd for OperationHub

[Path]
PathChanged=/etc/passwd
PathChanged=/etc/shadow
PathChanged=/etc/group
PathModified=/etc/passwd
PathModified=/etc/shadow
PathModified=/etc/group

[Install]
WantedBy=multi-user.target

